# scripts/ai.awfulclaw.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.awfulclaw.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/scripts/reload-on-main.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>WatchPaths</key>
    <array>
        <string>__PROJECT_DIR__/app/awfulclaw</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/watcher.out.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/watcher.err.log</string>
</dict>
</plist>
