# launchd/ai.awfulclaw.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.awfulclaw.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/scripts/start_agent.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/awfulclaw.out.log</string>
</dict>
</plist>
